# profile.d/bash.sh
FIGNORE=$FIGNORE:.reply:DS_Store:.hi:.pyc:.o
export FIGNORE
